# files/serv/condor-schedd
#!/bin/bash
#
# chkconfig: 35 90 12
# description: Condor startd
#
# Get function from functions library
#. /etc/init.d/functions
# Start the daemon
stat=0

start() {
        echo -n "Starting condor_schedd: "
        if [ -e /var/lock/condor_schedd ]
        then
          echo -n "condor_schedd already running"
        else
          touch /var/lock/condor_schedd
          /usr/sbin/condor_schedd &
          echo -n "condor_schedd started successfully"
        fi
        echo
}
# Stop the service
stop() {
        echo -n "Stopping condor_schedd: "
        pkill condor_schedd
        ### Now, delete the lock file ###
        rm -f /var/lock/condor_schedd
        echo -n "condor_schedd stopped successfully"
        echo
}
status() {
        pidof condor_schedd
        if [ $? -ne 0 ]
        then
          echo -n "condor_schedd is not running, "
          stat=1
        else
          echo -n "condor_schedd is running, "
        fi

        if [ -e /var/lock/condor_schedd ]
        then
          echo -n "lockfile present"
        else
          echo -n "lockfile not present"
        fi
        
        echo
}
### main logic ###
case "$1" in
  start)
        start
        ;;
  stop)
        stop
        ;;
  status)
        status
        ;;
  restart|reload|condrestart)
        stop
        start
        ;;
  *)
        echo $"Usage: $0 {start|stop|restart|reload|status}"
        exit 1
esac
exit $stat
